# supervisord/setup.conf
[program:setup-wp]
command=/bin/bash -c "/src/setup-wp.sh"